# systemd/desktop-entry-daemon.service
[Unit]
Description=XDG Temp daemon
PartOf=graphical-session.target

[Service]
BusName=net.ryanabx.DesktopEntry
Environment="RUST_LOG=debug"
ExecStart=desktop-entry-daemon
Type=dbus
Slice=session.slice

[Install]
WantedBy=graphical-session.target
Alias=dbus-net.ryanabx.DesktopEntry.service