# data/xdg-temp-daemon.service.in
[Unit]
Description=XDG Temp daemon
PartOf=graphical-session-pre.target
Before=graphical-session.target

[Service]
BusName=net.ryanabx.XDGTempDaemon
ExecStart=/usr/libexec/xdg-temp-daemon
Type=dbus
Slice=session.slice

[Install]
WantedBy=graphical-session-pre.target
Alias=dbus-net.ryanabx.XDGTempDaemon.service