# data/xdg-temp-daemon.service.in
[Unit]
Description=XDG Temp daemon
PartOf=graphical-session.target

[Service]
BusName=net.ryanabx.XDGTempDaemon
ExecStart=/usr/libexec/xdg-temp-daemon
Type=dbus
Slice=session.slice