# data/xdg-temp-daemon.service.in
[Unit]
Description=XDG Temp daemon
Before=multi-user.target display-manager.target

[Service]
Type=dbus
BusName=net.ryanabx.XDGTempDaemon
ExecStart=/usr/libexec/xdg-temp-daemon
Restart=on-failure

[Install]
WantedBy=graphical.target
